# packaging/agent/archive-agent.service
[Unit]
Description=comfy-chatbot archive agent (privileged zuluCrypt-cli mount helper)
Documentation=https://github.com/nerwander/comfy-chatbot
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/archive-agent
Restart=on-failure
RestartSec=2
# The agent runs as root: it mounts encrypted volumes via zuluCrypt-cli and sets
# up MOUNT_DIR as a shared mount. Mount propagation must reach other namespaces
# (the container), so the service shares the host mount namespace.
MountFlags=shared

[Install]
WantedBy=multi-user.target
