# packaging/deb/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Mount directories for the encrypted volumes (kept in sync with the
        # MOUNT_DIR / OUTPUT_MOUNT_DIR defaults in /etc/archive-agent.conf).
        mkdir -p /var/lib/archive-agent/mnt
        mkdir -p /var/lib/archive-agent/output

        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
            systemctl enable archive-agent.service || true
            systemctl restart archive-agent.service || true
        fi
        ;;
esac

exit 0
