# packaging/deb/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Mount directory for the encrypted volume (kept in sync with the
        # MOUNT_DIR default in /etc/archive-agent.conf).
        mkdir -p /var/lib/archive-agent/mnt

        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
            systemctl enable archive-agent.service || true
            systemctl restart archive-agent.service || true
        fi
        ;;
esac

exit 0
